# roles/ic880a_rpi_backplane/templates/env.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Run env.service every 30 secs

[Timer]
OnCalendar=*:*:5

[Install]
WantedBy=timers.target
