# roles/ic880a_rpi_backplane/templates/env.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Prints temperature and humidity into /tmp/temperature and /tmp/humidity files

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/bin/env.sh'
