# roles/ic880a_rpi_backplane/templates/voltage.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Prints voltage into /tmp/voltage file

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/bin/voltage.py > /tmp/voltage'

